# mylar/sv/mylar/finish
#!/bin/sh

kill -9 $(ps ax |grep Mylar.py |awk '{print $1}')
