# transmission/sv/transmission/finish
#!/bin/sh
## transmission.finish
## vers = 83bf920a131d262622091b1718ec8295

exec /bin/kill -s SIGKILL `cat /var/run/transmission.pid`
