# couchpotato/sv/couchpotato/finish
#!/bin/sh

kill -9 $(ps ax |grep CouchPotato.py | awk '{print $1}')
