# headphones/sv/headphones/finish
#!/bin/sh

kill -9 $(ps ax |grep Headphones.py |awk '{print $1}')
