# nzbget/sv/nzbget/finish
#!/bin/bash

/opt/nzbget/nzbget -c /config/nzbget.conf -Q
