# base/sv/firstboot/run
#!/bin/sh
## firstboot.run
## vers = 0ce2b7d9408756d6e2e4382f5b30f15f

if [ -f /deploy ]; then
	echo "[START] First Boot Deployment executing..."
	exec /opt/rootwyrm/deploy.sh
fi
