# sabnzbd/sv/sabnzbd/finish
#!/bin/bash

for pid in `ls /config/sabnzbd*.pid`; do
	kill -9 $pid
done
